# config/server/sidekiq.service.d/env.conf
[Service]
Environment="DISABLE_SPRING=1"
Environment="DEPLOYER=ec2-user"